# deploy/traektoria51-bot.service
[Unit]
Description=Traektoria51 bot (Telegram/MAX/Alice/web-chat)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
# Пользователь, от которого работает бот (на проде Beget — botsvc)
User=botsvc
Group=botsvc

WorkingDirectory=/opt/telegram-bot
# Путь к node уточните командой: which node  (на проде /usr/bin/node, v18)
ExecStart=/usr/bin/node /opt/telegram-bot/src/index.js

# Переменные берутся из .env (его читает сам код через src/load-env.js)
Environment=NODE_ENV=production

Restart=always
RestartSec=3
# Логи идут в journald: journalctl -u traektoria51-bot -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=traektoria51-bot

# Безопасность
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
